# src/kde/dot-config/systemd/user/plasma-powerdevil.service.d/no-ddcutil.conf
[Service]
# Disable DDC/CI brightness control for external monitors.
# Without this, powerdevil sends DDC commands via ddcutil to connected
# monitors over HDMI/DP, overriding your manual OSD brightness settings
# every time you unlock, wake from suspend, or change power profiles.
Environment=POWERDEVIL_NO_DDCUTIL=1
